# deploy/autodm.service
[Unit]
Description=Auto DM - Next.js production (autodm.pt)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/AutoDm
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -H 127.0.0.1 -p 3001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
